# kitchen/cookbooks/uwsgi/templates/uwsgi.service.erb
[Unit]
Description=uWSGI instance to serve <%= @project_name %>

[Service]
User=root
Group=root
ExecStart=/bin/bash -c 'uwsgi --ini /etc/uwsgi/apps-enabled/<%= @project_name %>.ini'

[Install]
WantedBy=multi-user.target