# kitchen/cookbooks/uwsgi/templates/uwsgi_service.conf.erb
description "uWSGI instance to serve <%= @project_name %>"

start on runlevel [2345]
stop on runlevel [!2345]

setuid vagrant
setgid vagrant

script
    cd <%= @virtual_env %>
    ./bin/activate
    uwsgi --ini .uwsgi/<%= @project_name %>.ini
end script
